# quantumd/assets/templates/systemd/quantumd.service.template
[Unit]
Description=Quantum Storage Daemon
After=network.target

[Service]
ExecStart=/usr/local/bin/quantumd
Restart=always
User=root

[Install]
WantedBy=multi-user.target
